# ews/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

# SSSD — auth AD via PAM (flag -i = interactive/foreground pour supervisord)
[program:sssd]
command=/usr/sbin/sssd -i --logger=stderr
autorestart=true
stdout_logfile=/var/log/sssd.log
stderr_logfile=/var/log/sssd.err
priority=5

# SSH avec auth PAM (→ SSSD → AD)
[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true
stdout_logfile=/var/log/sshd.log
stderr_logfile=/var/log/sshd.err
priority=8

# Affichage virtuel
[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x720x24 -ac +extension GLX +render -noreset
autorestart=true
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err
priority=10

# Bureau XFCE en tant qu'engineer (compte local)
[program:xfce]
command=su - engineer -c "export DISPLAY=:1; exec startxfce4"
autorestart=true
stdout_logfile=/var/log/xfce.log
stderr_logfile=/var/log/xfce.err
priority=30
startsecs=3

# VNC avec authentification par mot de passe (SSH pour auth AD)
[program:x11vnc]
command=/bin/bash -lc "x11vnc -display :1 -forever -shared -rfbport 5900 -rfbauth /etc/x11vnc.pass -noxdamage -o /var/log/x11vnc.log"
autorestart=true
stdout_logfile=/var/log/x11vnc.out
stderr_logfile=/var/log/x11vnc.err
priority=40
startsecs=5

# noVNC HTTPS avec certificat auto-signé
[program:novnc]
command=/bin/bash -lc "websockify --web /opt/noVNC 6080 localhost:5900 --cert /etc/ssl/novnc/novnc.pem"
autorestart=true
stdout_logfile=/var/log/novnc.out
stderr_logfile=/var/log/novnc.err
priority=50
startsecs=7
